# configs/server/wstunnel.service
[Unit]
Description=wstunnel WebSocket tunnel server
After=network.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wstunnel server ws://0.0.0.0:5555 --restrict-to 127.0.0.1:51820
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
